# Day4/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

reboot_ap_sta_detect.sh > /dev/null

exit 0
